# crontab
*   *   *   *   *   php /opt/app/weather/bin/console weather:predictions:collect 10 --env=prod
