# crontab
*   *   *   *   *   php /opt/app/bin/console weather:predictions:collect 10 --env=prod
